# deploy/systemd/xome-research-agent.service
[Unit]
Description=Xome Learning Brief research planner
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/Intern_Newsletter
ExecStart=/usr/bin/npm run research:plan
Environment=NODE_ENV=production
