# deploy/cron/xome-newsletter.cron
CRON_TZ=America/Chicago
0 7 * * * cd /root/Intern_Newsletter && ./scripts/run-daily-newsletter.sh
